# launchd/com.vacationbot.whatsapp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vacationbot.whatsapp</string>

    <key>UserName</key>
    <string>jon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jon/projects/vacation-bot/run_whatsapp.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jon/projects/vacation-bot</string>

    <!-- Start automatically at boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/Users/jon/projects/vacation-bot/logs/whatsapp.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jon/projects/vacation-bot/logs/whatsapp.error.log</string>
</dict>
</plist>
